# conf/etc/supervisord.ini
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

[program:v2ray]
command=/usr/local/bin/v2ray --config /etc/config.json
autorestart=true

[program:nginx]
command=/usr/sbin/nginx  -g 'daemon off;'
autorestart=true
